# deploy/coffee-web.service
[Unit]
Description=Arbor Web Report
After=network.target

[Service]
Type=simple
User=coffee
Group=coffee
WorkingDirectory=/home/coffee/coffee-v3
Environment=PYTHONPATH=/home/coffee/coffee-v3
Environment=PATH=/home/coffee/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3.11 -m uvicorn web.app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
